# systemd/weather_theme.timer
[Unit]
Description=Runs weather_theme every 30min

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
